# .config/systemd/user/rnsd.service
[Unit]
Description=Reticulum Network Stack
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/rnsd
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
